# ntopng/rootfs/etc/services.d/ntopng/run
#!/command/with-contenv bashio
# shellcheck shell=bash

echo "Starting ntopng..."
INTERFACE="$(bashio::config interface)"
echo "$INTERFACE"
sleep 10

exec ntopng \
  -i "${INTERFACE}" \
  -w 4000