# ntopng/rootfs/etc/services.d/ntopng/run
#!/bin/sh

echo "Starting ntopng..."
# INTERFACE="$(bashio::config interface)"
# echo "$INTERFACE"
# sleep 10
exec ntopng -i wlp1s0u1u4 -w 4000
# exec ntopng \
#   -i "${INTERFACE}" \
#   -w 4000