# ntopng/rootfs/etc/services.d/ntopng/run
#!/command/with-contenv bashio
# shellcheck shell=bash

echo "Starting ntopng..."
sleep infinity
# INTERFACE="$(bashio::config interface)"

# exec ntopng \
#   -i "${INTERFACE}" \
#   -w 4000