# grafana_stack/rootfs/etc/services.d/node-exporter/run
#!/command/with-contenv bashio
# shellcheck shell=bash

exec node_exporter