# ntopng/rootfs/etc/services.d/redis/run
#!/bin/sh
echo "start redis"
exec redis-server