# grafana_stack/rootfs/etc/services.d/prometheus/run
#!/command/with-contenv bashio
# shellcheck shell=bash

set -euo pipefail

PORT="$(bashio::config prometheus_port)"
RETENTION="$(bashio::config prometheus_retention_time)"

bashio::log.info "Starting Prometheus on port ${PORT}"

which prometheus
id
ls -ld /data

mkdir -p /data/prometheus

exec prometheus \
  --config.file=/data/prometheus/prometheus.yml \
  --storage.tsdb.path=/data/prometheus