# grafana_stack/rootfs/etc/services.d/iptable-export/run
#!/command/with-contenv bashio
# shellcheck shell=bash
while true; do

  TMP="/data/node_exporter/iptables.prom.tmp"

  iptables -L TRAFFIC_MONITOR -v -n -x | awk '
  /RETURN/ {
      bytes=$2
      ip=$8
      printf "host_upload_bytes_total{ip=\"%s\"} %s\n", ip, bytes
  }' > "$TMP"

  mv "$TMP" /data/node_exporter/iptables.prom

  sleep 5
done